# tests/RedisConfigs/Docker/supervisord.conf
[supervisord]
nodaemon=false

[program:master-6479]
command=/usr/local/bin/redis-server /data/Basic/master-6479.conf
directory=/data/Basic
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:replica-6480]
command=/usr/local/bin/redis-server /data/Basic/replica-6480.conf
directory=/data/Basic
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:secure-6481]
command=/usr/local/bin/redis-server /data/Basic/secure-6481.conf
directory=/data/Basic
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:master-6482]
command=/usr/local/bin/redis-server /data/Failover/master-6482.conf
directory=/data/Failover
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:replica-6483]
command=/usr/local/bin/redis-server /data/Failover/replica-6483.conf
directory=/data/Failover
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:cluster-7000]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7000.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:cluster-7001]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7001.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:cluster-7002]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7002.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:cluster-7003]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7003.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:cluster-7004]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7004.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:cluster-7005]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7005.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:redis-7010]
command=/usr/local/bin/redis-server /data/Sentinel/redis-7010.conf
directory=/data/Sentinel
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:redis-7011]
command=/usr/local/bin/redis-server /data/Sentinel/redis-7011.conf
directory=/data/Sentinel
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:sentinel-26479]
command=/usr/local/bin/redis-server /data/Sentinel/sentinel-26479.conf --sentinel
directory=/data/Sentinel
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:sentinel-26480]
command=/usr/local/bin/redis-server /data/Sentinel/sentinel-26480.conf --sentinel
directory=/data/Sentinel
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:sentinel-26481]
command=/usr/local/bin/redis-server /data/Sentinel/sentinel-26481.conf --sentinel
directory=/data/Sentinel
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true